# machines/nordeck-sip-template/etc/systemd/system/jitsi-component-sidecar-config.service
[Unit]
Description=Jitsi Component Sidecar Config
Wants=network-online.target
After=network-online.target
Before=jitsi-component-sidecar.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/jitsi-component-sidecar-config

[Install]
WantedBy=multi-user.target
